# cron/crontab
*/10 * * * * ~/matomepp/backend/bin/matomepp pickuptweet
