# cron/crontab
*/5 * * * * ~/matomepp/backend/bin/matomepp pickuptweet --newest
*/10 * * * * ~/matomepp/backend/bin/matomepp pickuptweet --random
