# cron/crontab
0 0 * * * ~/matomepp/backend/bin/matomepp pickuptweet $((RANDOM % 3000)) 2> /dev/null
